# dot_config/autostart/evolution.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=evolution
Comment=evolution
Exec=evolution
StartupNotify=false
Terminal=true
